# templates/default/sonarr_upstart.conf.erb
author "Simon Tallmyr - Nosscire"
description "Upstart Script to run Sonarr as a service"

setuid <%= @user %>

env DIR=<%= @daemon_path %>

setgid nogroup
start on runlevel [2345]
stop on runlevel [016]

respawn

exec mono --debug $DIR/NzbDrone.exe
