# templates/default/sonarr_systemd.service.erb
[Unit]
Description=Sonarr Daemon
After=syslog.target network.target

[Service]
User=<%= @service_username %>
Group=<%= @service_group %>

Type=simple
ExecStart=<%= @mono_path %> <%= @sonarr_exe %> -nobrowser
TimeoutStopSec=20
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
